# puppet_code_implementation/systemd_unit_file/vlc.service
[Unit]
Description=vlc.service

[Service]
User=root
Environment="DISPLAY=:0"
WorkingDirectory=/
Restart=always

[Install]
WantedBy=default.target